# clients/supervisor/conf.d/rtorrent.supervisor.conf
[program:rtorrent]
command=rtorrent
user=toto
autostart=true
autorestart=true
stdout_logfile=/var/rtorrent/logs/rtorrent_out.log
stderr_logfile=/var/rtorrent/logs/rtorrent_err.log